# tests/scenarios/scenario_a_basic_multi_proc_host/host-sim/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:api]
command=/usr/bin/python3 /opt/api/app.py
directory=/opt/api
user=api-user
autostart=true
autorestart=true
stdout_logfile=/var/log/apps/api.log
stderr_logfile=/var/log/apps/api.err
environment=
    DATABASE_URL="postgresql://app:secret123@db:5432/appdb",
    REDIS_URL="redis://redis:6379",
    API_PORT="8080",
    LOG_LEVEL="info"

[program:worker]
command=/usr/bin/python3 /opt/worker/worker.py
directory=/opt/worker
user=worker-user
autostart=true
autorestart=true
stdout_logfile=/var/log/apps/worker.log
stderr_logfile=/var/log/apps/worker.err
environment=
    DATABASE_URL="postgresql://app:secret123@db:5432/appdb",
    REDIS_URL="redis://redis:6379",
    WORKER_CONCURRENCY="4"

[program:metrics-exporter]
command=/usr/bin/python3 /opt/metrics/exporter.py
directory=/opt/metrics
autostart=true
autorestart=true
stdout_logfile=/var/log/apps/metrics.log
stderr_logfile=/var/log/apps/metrics.err
environment=
    METRICS_PORT="8081"
